# infra/systemd/cynic-healthcheck.service
[Unit]
Description=CYNIC Readiness Check
After=cynic-kernel.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/cynic/env
# Rule #15: check HTTP status code, don't parse JSON. /ready returns 200 or 503.
ExecStart=/usr/bin/curl -sf --max-time 10 http://${CYNIC_REST_ADDR}/ready
