# infra/systemd/cynic-spike-detector.service
[Unit]
Description=CYNIC Spike Detector — GeckoTerminal trending pools → /observe domain=token-analysis
After=cynic-kernel.service
Wants=network-online.target
After=network-online.target

[Service]
EnvironmentFile=%h/.config/cynic/env
Type=oneshot
WorkingDirectory=%h/cynic/cynic-python/sensors
ExecStart=/usr/bin/python3 %h/cynic/cynic-python/sensors/spike_detector.py
TimeoutStartSec=60
TimeoutStopSec=10
MemoryMax=128M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cynic-spike-detector

[Install]
WantedBy=timers.target
