# infra/systemd/cynic-node.service
[Unit]
Description=CYNIC Node — WebSocket Federation Client
After=network.target
Wants=cynic-kernel.service

[Service]
Type=simple
WorkingDirectory=%h/Bureau/CYNIC
# Generate federation endpoints + resolve node config placeholders at startup
ExecStartPre=%h/Bureau/CYNIC/scripts/fleet-gen.py --mode node-kernels
ExecStartPre=/bin/bash -c '%h/Bureau/CYNIC/scripts/fleet-gen.py %h/Bureau/CYNIC/cynic-node/node-qwen35-gpu.toml > %h/.config/cynic/cynic-node.toml'
ExecStart=%h/bin/cynic-node --config %h/.config/cynic/cynic-node.toml
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5
# Environment from ~/.cynic-env (Bearer token, REST addr, API key)
EnvironmentFile=%h/.config/cynic/env
Environment="CYNIC_SCRIPTS=%h/Bureau/CYNIC/scripts"
LimitNOFILE=65535
MemoryMax=1G

[Install]
WantedBy=default.target
