# infra/systemd/cynic-monitor.service
[Unit]
Description=CYNIC Proxmox Monitoring Sonde

[Service]
Type=oneshot
ExecStart=/bin/bash %h/cynic/infra/monitoring/monitor_vms.sh
WorkingDirectory=%h/cynic

