# infra/systemd/cynic-healthcheck.service
[Unit]
Description=CYNIC Health Check

[Service]
Type=oneshot
ExecStart=%h/bin/cynic-healthcheck.sh
