# infra/systemd/cynic-healthcheck.service
[Unit]
Description=CYNIC Health Check — Telegram alerting on kernel degradation
After=cynic-kernel.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/cynic/env
ExecStart=%h/Bureau/CYNIC/scripts/cynic-healthcheck.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cynic-healthcheck
