# infra/systemd/cynic-monitor.timer
[Unit]
Description=Run CYNIC Proxmox Monitoring every minute

[Timer]
OnBootSec=5min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
