# infra/systemd/hermes-search-executor.service
[Unit]
Description=CYNIC Hermes Search Executor — automated X.com searches
After=cynic-kernel.service
PartOf=hermes-organism.target

[Service]
EnvironmentFile=%h/.config/cynic/env
Type=oneshot
Environment=X_ORGAN_DIR=%h/.cynic/organs/hermes/x
Environment=X_PROXY_ADDR=localhost:8888
WorkingDirectory=%h/Bureau/CYNIC/scripts/hermes-x
# CYNIC_REST_ADDR is injected from EnvironmentFile (should be substituted by config-sync.sh hermes-config)
ExecStart=/usr/bin/python3 core/search_executor.py \
    --organ-dir %h/.cynic/organs/hermes/x \
    --max-searches 5 \
    --proxy-addr localhost:8888
TimeoutStartSec=600
TimeoutStopSec=30
MemoryMax=512M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-search-executor

[Install]
WantedBy=timers.target
