# infra/systemd/hermes-x-ingest.service
[Unit]
Description=CYNIC Hermes X ingest — tails dataset.jsonl → POST /observe
After=cynic-kernel.service hermes-proxy.service

[Service]
Type=simple
Environment=X_DATASET_PATH=%h/.cynic/organs/hermes/x/dataset.jsonl
EnvironmentFile=-%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/scripts/hermes-x/core/x_ingest_daemon.py
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
