# infra/systemd/hermes-k15-consumer.service
[Unit]
Description=K15 Consumer — Score observations and dispatch high-signal tasks
Requires=network.target
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/Bureau/CYNIC

# Run the K15 consumer via wrapper script that sources ~/.cynic-env
ExecStart=/home/user/Bureau/CYNIC/scripts/start-k15-consumer.sh

# Restart on failure with 30s backoff
Restart=on-failure
RestartSec=30

# Output to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
