# infra/systemd/hermes-x-ingest.service
[Unit]
Description=CYNIC Hermes X ingest — tails dataset.jsonl → POST /observe
After=cynic-kernel.service hermes-proxy.service

[Service]
Type=simple
Environment=X_DATASET_PATH=%h/.cynic/organs/hermes/x/dataset.jsonl
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-%h/.config/cynic/env
ExecStartPre=/bin/sh -c 'test -n "$CYNIC_REST_ADDR" || { echo "CYNIC_REST_ADDR not set — check ~/.config/cynic/env"; exit 1; }'
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/scripts/hermes-x/core/x_ingest_daemon.py
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
