# infra/systemd/hermes-k15-consumer.service
[Unit]
Description=CYNIC K15 Observation Consumer — Score observations, dispatch hermes tasks
After=cynic-kernel.service
PartOf=hermes-k15-consumer.timer

[Service]
Type=oneshot
Environment="HOME=/home/user"
EnvironmentFile=%h/.config/cynic/env
EnvironmentFile=%h/.cynic-env
WorkingDirectory=%h/Bureau/CYNIC
ExecStart=/usr/bin/python3 cynic-python/consumers/k15_observation_consumer.py --kernel-url http://<TAILSCALE_CORE>:3030 --poll-interval 10 --max-iter 1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-k15-consumer

[Install]
WantedBy=hermes-k15-consumer.timer
