# infra/systemd/hermes-search-executor.service
[Unit]
Description=Hermes X Search Executor — stochastic behavioral search
After=hermes-browser.service
Wants=hermes-browser.service
PartOf=hermes-search-executor.timer

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/cynic/env
WorkingDirectory=/home/user/Bureau/CYNIC/scripts/hermes-x
ExecStart=/usr/bin/python3 search_executor.py \
    --organ-dir /home/user/.cynic/organs/hermes/x \
    --max-searches 5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-search-executor
MemoryMax=512M
TimeoutStartSec=120
