# infra/systemd/hermes-cost-aggregator.timer
[Unit]
Description=CYNIC Metabolic Cost Aggregator timer — every 30min

[Timer]
OnCalendar=*:00/30
Persistent=true

[Install]
WantedBy=timers.target
