# infra/systemd/hermes-task-runner.service
[Unit]
Description=CYNIC Hermes task runner — claims organ-local tasks and executes
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=%h/.config/cynic/env
Environment="X_ORGAN_DIR=%h/.cynic/organs/hermes/x"
ExecStart=/usr/bin/python3 %h/cynic/scripts/hermes-x/core/hermes_task_runner.py --organ-dir $X_ORGAN_DIR --interval 10
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
