# infra/systemd/hermes-browser@.service
[Unit]
Description=CYNIC Hermes browser — %i account — Chrome with proxy + CDP for passive capture
# Hub depends on us, not the other way around. Proxy is independent.
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=%h/.config/cynic/env
Environment=CHROME_BIN=/usr/bin/google-chrome
Environment=HERMES_ACCOUNT=%i
Environment=HERMES_CDP_PORT=40769
Environment=HERMES_PROXY_PORT=8888
Environment=DISPLAY=:0
# Account-specific profile (loaded by launch-browser.sh from accounts.toml)
Environment=HERMES_CHROME_PROFILE=%h/.cynic/organs/hermes/x/chrome-profiles/%i
ExecStart=%h/Bureau/CYNIC/scripts/hermes-x/scripts/launch-browser.sh
Restart=on-failure
RestartSec=30
TimeoutStopSec=10

[Install]
WantedBy=default.target
