# infra/systemd/hermes-agent-executor.service
[Unit]
Description=CYNIC Hermes Agent Executor — autonomous meta-agent task runner
After=cynic-kernel.service hermes-x-ingest.service
Wants=hermes-x-ingest.service

[Service]
Type=simple
EnvironmentFile=%h/.config/cynic/env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/cynic/services/cynic-python
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HERMES_AGENT_DOMAIN_DENY=organ-anvil,repo-lifecycle
ExecStart=/usr/bin/python3 -m organs.core.daemon --organ-dir %h/.cynic/organs/hermes/x
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
