# infra/systemd/hermes-gemini-briefing.service
[Unit]
Description=CYNIC Hermes Gemini Briefing — synthesizes domain wisdom (4h cron)
After=cynic-kernel.service hermes-x-ingest.service
PartOf=hermes-gemini-briefing.timer

[Service]
Type=oneshot
EnvironmentFile=%h/.config/cynic/env
Environment=X_ORGAN_DIR=%h/.cynic/organs/hermes/x
ExecStart=/usr/bin/python3 %h/cynic/scripts/hermes-x/analysis/gemini_briefing_consumer.py --organ-dir %h/.cynic/organs/hermes/x
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=hermes-gemini-briefing.timer
