# infra/systemd/hermes-search-executor.timer
[Unit]
Description=Hermes X Search Executor — stochastic timing (wide window)
Requires=hermes-search-executor.service

[Timer]
OnCalendar=*-*-* *:05,20,35,50:00
RandomizedDelaySec=10min
Persistent=false

[Install]
WantedBy=timers.target
