# infra/systemd/hermes-x-login@.service
[Unit]
Description=CYNIC Hermes X.com authentication — %i account — login via hermes-browser CDP
After=hermes-browser@%i.service
Requires=hermes-browser@%i.service

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/cynic/env
Environment=HERMES_ACCOUNT=%i
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/user/Bureau/CYNIC/scripts/hermes-x/core/hermes_x_login.py
RemainAfterExit=yes
TimeoutStartSec=60

[Install]
WantedBy=default.target
