# infra/systemd/hermes-proxy.service
[Unit]
Description=CYNIC Hermes X proxy (mitmdump)
After=cynic-kernel.service

[Service]
Type=simple
# R1: project root from env, not hardcoded path
EnvironmentFile=%h/.config/cynic/env
WorkingDirectory=${CYNIC_PROJECT_ROOT}/scripts/hermes-x
ExecStart=%h/.local/bin/mitmdump -s x_proxy.py --listen-host 127.0.0.1 -p 8888 --set stream_large_bodies=1m --set ssl_insecure=true -q
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
