# infra/systemd/hermes-browser.service
[Unit]
Description=CYNIC Hermes browser — Chrome with proxy + CDP for passive capture and agent control
# Hub depends on us, not the other way around. Proxy is independent.
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=%h/.config/cynic/env
Environment=CHROME_BIN=/usr/bin/google-chrome
Environment=HERMES_CDP_PORT=40769
Environment=HERMES_PROXY_PORT=8888
Environment=DISPLAY=:0
ExecStart=%h/Bureau/CYNIC/scripts/hermes-x/scripts/launch-browser.sh
Restart=on-failure
RestartSec=30
TimeoutStopSec=10

[Install]
WantedBy=default.target
