# infra/systemd/hermes-proxy.service
[Unit]
Description=CYNIC Hermes proxy — passive X/Twitter capture (mitmdump)
After=cynic-kernel.service

[Service]
Type=simple
# R1: project root from env, not hardcoded path
EnvironmentFile=%h/.config/cynic/env
WorkingDirectory=%h/Bureau/CYNIC/scripts/hermes-x
Environment=X_CAPTURE_DIR=%h/.cynic/organs/hermes/x/captures
Environment=X_DATASET_PATH=%h/.cynic/organs/hermes/x/dataset.jsonl
# SoC: only decrypt X/Twitter traffic. All other HTTPS passes through opaque.
ExecStart=%h/.local/bin/mitmdump \
    -s x_proxy.py \
    --listen-host 127.0.0.1 -p 8888 \
    --set stream_large_bodies=1m \
    --set ssl_insecure=true \
    --set allow_hosts="(?:.*\\.)?x\\.com|(?:.*\\.)?twitter\\.com" \
    -q
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=default.target
