# infra/systemd/hermes-x-recovery@.service
[Unit]
Description=CYNIC Hermes X recovery daemon — %i account
After=network-online.target hermes-x-ingest@%i.service
Wants=network-online.target
PartOf=hermes-x-ingest@%i.service

[Service]
Type=simple
EnvironmentFile=%h/.config/cynic/env
Environment=HERMES_ACCOUNT=%i
Environment=RECOVERY_LOG_PATH=%h/.cynic/organs/hermes/x/recovery.log
Environment=RECOVERY_CHECK_INTERVAL=30
WorkingDirectory=%h/Bureau/CYNIC/scripts/hermes-x/core
ExecStart=/usr/bin/python3 hermes_x_recovery.py
Restart=on-failure
RestartSec=10
# Recovery daemon must keep running to detect failures — aggressive restart
StandardOutput=append:%h/.cynic/organs/hermes/x/recovery.log
StandardError=append:%h/.cynic/organs/hermes/x/recovery.log
# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
