# infra/systemd/mirror-agent.service
[Unit]
Description=CYNIC Mirror Agent v0.1.0 — behavioral model daemon
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.cynic/organs/mirror/.venv/bin/python -m organs.mirror.daemon
WorkingDirectory=%h/Bureau/CYNIC/cynic-python
EnvironmentFile=%h/.config/cynic/env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
