# infra/systemd/organ-anvil-hermes-agent.service
[Unit]
Description=CYNIC Organ Anvil Hermes Agent - repo lifecycle manager
After=cynic-kernel.service hermes-agent-executor.service
Wants=cynic-kernel.service

[Service]
Type=simple
WorkingDirectory=%h/Bureau/CYNIC
EnvironmentFile=%h/.config/cynic/env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=HERMES_AGENT_DOMAIN_ALLOW=organ-anvil,repo-lifecycle
Environment=HERMES_AGENT_KIND_ALLOW=hermes
Environment=HERMES_AGENT_POLL_LIMIT=100
Environment=X_ORGAN_DIR=%h/.cynic/organs/anvil
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/scripts/hermes-x/core/hermes_agent_task_executor.py --organ-dir %h/.cynic/organs/anvil --interval 30
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
