# infra/systemd/organ-anvil-cron.service
[Unit]
Description=CYNIC Organ Anvil Cron - dispatch repo lifecycle task
After=cynic-kernel.service organ-anvil-hermes-agent.service
Wants=organ-anvil-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/Bureau/CYNIC
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/scripts/organ-anvil-dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
