# infra/systemd/organ-x-cron.service
[Unit]
Description=CYNIC Organ X Cron - dispatch X surface task
After=cynic-kernel.service organ-x-hermes-agent.service
Wants=organ-x-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/Bureau/CYNIC
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/scripts/organ-x-dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
