# infra/systemd/surrealdb-backup.timer
[Unit]
Description=Daily SurrealDB backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
