# infra/systemd/tailscale-fleet-probe.timer
[Unit]
Description=Tailscale Fleet Introspection Probe Timer
Documentation=docs/fleet-probe.md

[Timer]
# Run every 30 seconds
OnBootSec=5s
OnUnitActiveSec=30s

# Fire-and-forget: continue even if service fails
Unit=tailscale-fleet-probe.service

[Install]
WantedBy=timers.target
