# infra/systemd/telegram-listener.service
[Unit]
Description=CYNIC Telegram Listener
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Bureau/CYNIC/cynic-python
ExecStart=%h/.cynic/organs/telegram/.venv/bin/python -m organs.telegram.listener
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.config/cynic/env
EnvironmentFile=%h/.cynic/organs/telegram/env
# Cgroup limits — observed: 268MB.
MemoryMax=512M

[Install]
WantedBy=default.target
