# infra/systemd/token-snapshot.service
[Unit]
Description=Token Snapshot — daily holder + market + trajectory + outcomes + unify
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/cynic/env
WorkingDirectory=%h/Bureau/CYNIC
ExecStart=/usr/bin/python3 cynic-python/heuristics/collection/daily_snapshot.py
ExecStartPost=/usr/bin/python3 cynic-python/heuristics/collection/daily_snapshot.py --market
ExecStartPost=/usr/bin/python3 cynic-python/heuristics/collection/daily_snapshot.py --trajectory
ExecStartPost=/usr/bin/python3 cynic-python/heuristics/collection/outcome_collector.py
ExecStartPost=/usr/bin/python3 cynic-python/heuristics/collection/unify.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=token-snapshot
MemoryMax=256M
TimeoutStartSec=900

[Install]
WantedBy=token-snapshot.timer
