# organs/twitter/systemd/organ-twitter-cron.timer
[Unit]
Description=CYNIC Organ Twitter Cron Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
