# organs/vercel/systemd/organ-vercel-cron.service
[Unit]
Description=CYNIC Organ Vercel Cron - dispatch deployment surface task
After=cynic-kernel.service organ-vercel-hermes-agent.service
Wants=organ-vercel-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/cynic/organs/vercel/src/monitor.py
ExecStart=/usr/bin/python3 %h/cynic/organs/vercel/src/dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
