# organs/jupyter/systemd/organ-jupyter-cron.service
[Unit]
Description=CYNIC Organ - Jupyter dispatch observer

[Service]
Type=oneshot
WorkingDirectory=/home/user/Bureau/CYNIC
ExecStart=/home/user/Bureau/CYNIC/organs/jupyter/.venv/bin/python /home/user/Bureau/CYNIC/organs/jupyter/src/dispatch.py
