# organs/jupyter/systemd/organ-jupyter-cron.timer
[Unit]
Description=Timer for CYNIC Organ Jupyter Dispatch

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
