# organs/jupyter/systemd/organ-jupyter-lab.service
[Unit]
Description=CYNIC Organ - JupyterLab Server
After=network.target

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/Bureau/CYNIC
Environment="PATH=/home/user/Bureau/CYNIC/organs/jupyter/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="JUPYTER_PORT=8888"
Environment="JUPYTER_TOKEN=cynic-internal-token-2026"

# Lancement de JupyterLab sans navigateur, lié en localhost uniquement, avec le token fixé pour l'API
ExecStart=/home/user/Bureau/CYNIC/organs/jupyter/.venv/bin/jupyter-lab \
  --no-browser \
  --ip=127.0.0.1 \
  --port=8888 \
  --ServerApp.token="cynic-internal-token-2026" \
  --ServerApp.password="" \
  --ServerApp.root_dir=/home/user/Bureau/CYNIC

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
