# organs/jupyter/systemd/organ-jupyter-cron.service
[Unit]
Description=CYNIC Organ - Jupyter dispatch observer

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
ExecStart=%h/cynic/organs/jupyter/.venv/bin/python %h/cynic/organs/jupyter/src/server.py
