# organs/telegram/systemd/organ-telegram-cron.service
[Unit]
Description=CYNIC Organ Telegram Cron - dispatch messaging surface task
After=cynic-kernel.service organ-telegram-hermes-agent.service
Wants=organ-telegram-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/cynic/organs/telegram/src/dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
