# organs/forge/systemd/organ-forge-cron.service
[Unit]
Description=CYNIC Organ Anvil Cron - dispatch repo lifecycle task
After=cynic-kernel.service organ-forge-hermes-agent.service
Wants=organ-forge-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/Bureau/CYNIC
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/organs/forge/src/entrypoints/cli.py
ExecStart=/usr/bin/python3 %h/Bureau/CYNIC/organs/forge/src/dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
