# organs/forge/systemd/organ-forge-hermes-agent.service
[Unit]
Description=CYNIC Organ Anvil Hermes Agent - repo lifecycle manager
After=cynic-kernel.service hermes-agent-executor.service
Wants=cynic-kernel.service

[Service]
Type=simple
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/cynic/services/cynic-python
Environment=HERMES_AGENT_DOMAIN_ALLOW=organ-forge,repo-lifecycle
Environment=HERMES_AGENT_KIND_ALLOW=hermes
Environment=HERMES_AGENT_POLL_LIMIT=100
Environment=X_ORGAN_DIR=%h/cynic/organs/forge/state
ExecStart=/usr/bin/python3 -m organs.core.daemon --organ-dir %h/cynic/organs/forge/state --interval 30
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
