# organs/forge/systemd/organ-forge-cron.timer
[Unit]
Description=CYNIC Organ Anvil Cron - every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true
Unit=organ-forge-cron.service

[Install]
WantedBy=timers.target
