# scripts/cynic-ouroboros.timer
[Unit]
Description=Run CYNIC Ouroboros Loop nightly at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
