# organs/kanban-ui/systemd/organ-kanban-ui-cron.service
[Unit]
Description=CYNIC Organ Kanban UI Projector
After=cynic-kernel.service

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/cynic/organs/kanban-ui/src/entrypoints/cli.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
