# organs/kanban-ui/systemd/organ-kanban-ui-cron.timer
[Unit]
Description=CYNIC Organ Kanban UI Projector Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
