# organs/docker/systemd/organ-docker-hermes-agent.service
[Unit]
Description=CYNIC Organ Docker Hermes Agent - local container surface manager
After=cynic-kernel.service hermes-agent-executor.service docker.service
Wants=cynic-kernel.service docker.service

[Service]
Type=simple
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/cynic/services/cynic-python
Environment=HERMES_AGENT_DOMAIN_ALLOW=docker
Environment=HERMES_AGENT_KIND_ALLOW=hermes
Environment=HERMES_AGENT_POLL_LIMIT=100
Environment=X_ORGAN_DIR=%h/cynic/organs/docker/state
ExecStartPre=/usr/bin/install -Dm644 %h/cynic/infra/organ-docker/SKILL.md %h/.cynic/organs/docker/SKILL.md
ExecStart=/usr/bin/python3 -m organs.core.daemon --organ-dir %h/cynic/organs/docker/state --interval 30
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
