# organs/docker/systemd/organ-docker-cron.service
[Unit]
Description=CYNIC Organ Docker Cron - dispatch portal container task
After=cynic-kernel.service organ-docker-hermes-agent.service
Wants=organ-docker-hermes-agent.service

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/cynic/organs/docker/src/monitor.py
ExecStart=/usr/bin/python3 %h/cynic/organs/docker/src/dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
