# organs/reflection/systemd/organ-reflection-cron.service
[Unit]
Description=CYNIC Organ Reflection Cron
After=cynic-kernel.service

[Service]
Type=oneshot
WorkingDirectory=%h/cynic
EnvironmentFile=%h/.config/cynic/env
ExecStart=/usr/bin/python3 %h/cynic/organs/reflection/src/entrypoints/cli.py
ExecStart=/usr/bin/python3 %h/cynic/organs/reflection/src/dispatch.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
