# organs/reflection/systemd/organ-reflection-cron.timer
[Unit]
Description=CYNIC Organ Reflection Cron Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
